# src/bootloader/plymouth-quit.service.d/override.conf
[Unit]
Description=DELAYED termination of Plymouth Boot Screen
Before=lightdm.service

[Service]
ExecStart=
ExecStart=/bin/sleep 18
ExecStartPost=-/bin/plymouth quit
Type=forking
TimeoutSec=25

